# etc/supervisor/conf.d/control_mocp.conf
[program:control_mocp]
command=/root/control_mocp.sh
numprocs=1
stdout_logfile=/var/log/control_mocp.log
autostart=true
autorestart=true
user=root
stopsignal=KILL
